# docker/app/supervisor/conf.d/php-fpm.conf
[program:php]
process_name=%(program_name)s_%(process_num)02d
command=docker-php-entrypoint php-fpm # This is our main process so let us put it here
autostart=true
autorestart=true